# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:celeryd]
; set to your current virtualenv
user=romanstatkevich
command=/home/romanstatkevich/development/pets/birthday_bot/.env/bin/celery -A bot worker --loglevel=INFO
directory=/home/romanstatkevich/development/pets/birthday_bot/
stdout_logfile=/var/log/supervisor/celeryd.log
stderr_logfile=/var/log/supervisor/celeryd.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600

[program:celerybeat]
; set to your current virtualenv
command=/home/romanstatkevich/development/pets/birthday_bot/.env/bin/celery -A bot beat --loglevel=INFO
directory=/home/romanstatkevich/development/pets/birthday_bot/
stdout_logfile=/var/log/supervisor/celeryd_beat.log
stderr_logfile=/var/log/supervisor/celeryd_beat.log
user=romanstatkevich
numprocs=1
autostart=true
autorestart=true
startsecs=10

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

